# systemd/danzee-homelab-backup.service.d/uptime-kuma-heartbeat.conf
[Service]
ExecStartPost=/usr/local/sbin/danzee-backup-heartbeat
